# web/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:httpd]
command=apachectl -DFOREGROUND
user=root

[program:tomcat]
command=/usr/local/tomcat/bin/catalina.sh run
user=root
